# ansible/roles/wyoming-whisper/templates/wyoming-whisper.container.j2
[Unit]
Description=Wyoming Faster Whisper (GPU-accelerated)
After=network-online.target nvidia-persistenced.service
Wants=network-online.target

[Container]
Image=docker.io/rhasspy/wyoming-whisper:{{ wyoming_whisper_version | default('latest') }}
ContainerName=wyoming-whisper
PublishPort={{ wyoming_whisper_port | default('10300') }}:10300
Volume=/var/lib/wyoming-whisper:/data
AddDevice=nvidia.com/gpu=all
Environment=NVIDIA_VISIBLE_DEVICES=all
Environment=NVIDIA_DRIVER_CAPABILITIES=compute,utility
Exec=--model {{ wyoming_whisper_model | default('medium') }} --language {{ wyoming_whisper_language | default('en') }} --uri tcp://0.0.0.0:10300 --data-dir /data --download-dir /data

[Service]
Restart=always
RestartSec=10
# Model download and loading can take several minutes on first start
TimeoutStartSec=600

[Install]
WantedBy=multi-user.target
